# etc/crontab
*/5 * * * *	/bin/sh /app/data.sh
8 6-16 * * 2,3,4	/bin/sh /app/results.sh
